# meta-shatrox-rpi/recipes-core/startup-sound-service/files/startup-sound.service
[Unit]
Description=Startup Sound Service
After=sound.target systemd-user-sessions.service
Wants=sound.target

[Service]
Type=oneshot
ExecStart=/usr/bin/startup-sound.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
