# meta-shatrox-rpi/recipes-core/shatrox-motor-control/files/shatrox-motor-control.service
[Unit]
Description=SHATROX Motor Control Service
Documentation=https://github.com/shatrix/rpi5-rpios-ai-robot
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/bin/motor_controller.py
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shatrox-motor

# Security (relaxed for GPIO/I2C access)
PrivateTmp=false

# Resource limits
MemoryLimit=256M

[Install]
WantedBy=multi-user.target
